# templates/etc/systemd/system/check-mk-agent.socket.d/allowip.conf
[Socket]
IPAddressAllow=127.0.0.1 {% for i in lihas_mon_hosts | default([]) %} {{ i }} {% endfor %}

